# deploy/webterm.service
# systemd unit for webterm. Copy to /etc/systemd/system/webterm.service,
# edit User/paths/token, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now webterm
#
# NOTE: all terminal sessions live inside this process. Restarting this service
# kills every running session — that is the accepted trade-off of not using
# dtach/tmux. Avoid unnecessary restarts.

[Unit]
Description=webterm - persistent web terminal over HTTP/1.1
After=network.target

[Service]
Type=simple
User=youruser
WorkingDirectory=/opt/webterm
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=2

# Configuration (or use an EnvironmentFile=/opt/webterm/.env).
Environment=WEBTERM_HOST=127.0.0.1
Environment=WEBTERM_PORT=8080
# REQUIRED: set a strong random token before enabling the service. The server
# refuses to start if this is empty or a placeholder. Generate one with:
#   openssl rand -base64 32
# Prefer keeping the secret out of the unit file via EnvironmentFile below.
Environment=WEBTERM_TOKEN=
# Environment=WEBTERM_CMD=/bin/bash
# Environment=WEBTERM_ARGS=-l
# EnvironmentFile=/opt/webterm/.env

[Install]
WantedBy=multi-user.target
